# docker/t-api/supervisor/conf.d/crond.ini
[program:crond]
directory = /var/www/php
command = /usr/sbin/crond -n
autostart = true
autorestart = true
redirect_stderr = true
