# docker/t-api/supervisor/conf.d/laravel.default.ini
[program:laravel.default]
directory = /var/www/php
command = php artisan queue:work
autostart = true
autorestart = true
redirect_stderr = true
user = www-data

