# docker/t-api/supervisor/conf.d/laravel.sync.ini
[program:laravel.sync]
directory = /var/www/php
command = php artisan queue:work --queue=sync
autostart = true
autorestart = true
redirect_stderr = true
user = www-data
numprocs = 1
process_name = %(process_num)s
