# docker/t-api/supervisor/conf.d/laravel.crm.so.sync.ini
[program:laravel.crm]
directory = /var/www/php
command = php artisan crm-shiptask-queue
autostart = true
autorestart = true
redirect_stderr = true
user = www-data
numprocs = 1
process_name = %(process_num)s
