# docker/t-api/supervisor/conf.d/php-fpm.ini
[program:php-fpm]
directory = /var/www/php
command = php-fpm -O
autostart = true
autorestart = true
redirect_stderr = true
